# arboldes/usr/ports/mystuff/sysutils/smartmontools/pkg/smarts.rc
#!/bin/ksh

servicio="${TRUEPREFIX}/sbin/smartd"

. /etc/rc.d/rc.subr

rc_cmd $1
